# root/etc/systemd/system/sympa.service.d/nethserver.conf
[Unit]
# nethserver directive, sympa fail if we restart after mariadb
After=postfix.service
After=mariadb.service
After=mysqld.service
